# root/etc/s6-overlay/s6-rc.d/init-branding/up
#!/bin/sh
echo "Branding configured successfully"